# stacks/cortex-incus/cortex-incus-zpool.service.in
[Unit]
Description=Import CortexOS Incus file-backed ZFS pool
After=local-fs.target
Before=incus.service
RequiresMountsFor=@@POOL_SOURCE_DIR@@

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/sh -lc '/usr/sbin/zpool list @@POOL_NAME@@ >/dev/null 2>&1 || /usr/sbin/zpool import -N -d @@POOL_SOURCE_DIR@@ @@POOL_NAME@@'

[Install]
WantedBy=multi-user.target
