# packages/cortex-agent-generator/cortex-agent-generator.service
[Unit]
Description=CortexOS Agent Generator WS sidecar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# node-pty + sudo shells: run as root (mirrors cortex-terminal). The WS auth
# gate enforces admin + live OS group membership; we trust the gate, not the
# filesystem.
EnvironmentFile=/opt/cortexos/.secrets/dashboard.env
Environment=GENERATOR_PORT=3082
Environment=ALLOWED_ORIGIN=same-origin

WorkingDirectory=/opt/cortexos/packages/cortex-agent-generator
ExecStart=/usr/bin/env node ./src/server.js

Restart=always
RestartSec=2
RestartPreventExitStatus=0 2 3 4 5

# Hardening (mirrors cortex-terminal).
# Generator needs to write to /opt/cortexos/{hermes,.secrets,bin,templates}
# and /etc/systemd/system (rendered units), so ProtectSystem=full + explicit
# ReadWritePaths (not strict, which would make ReadWritePaths ignored).
ProtectSystem=full
ProtectHome=read-only
PrivateDevices=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
DevicePolicy=closed
PrivateUsers=no

ReadWritePaths=/opt/cortexos/hermes /opt/cortexos/.secrets /opt/cortexos/bin /opt/cortexos/templates /etc/systemd/system

[Install]
WantedBy=multi-user.target
