# packages/cortex-agent-generator/cortex-agent-generator.service
[Unit]
Description=CortexOS Agent Generator WS sidecar
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

# Hardening mirrors cortex-terminal.service: node-pty must fork shells, so
# this is deliberately looser than the dashboard unit (no ProtectSystem, no
# PrivateUsers — needs /dev/pts; no PrivateDevices — needs /dev/ptmx; no
# MemoryDenyWriteExecute — node-pty allocates rxvt-unicode-256color; no
# RestrictNamespaces — node-pty may need CLONE_NEWNS for the slave PTY).
EnvironmentFile=/opt/cortexos/.secrets/dashboard.env
Environment=GENERATOR_PORT=3082
Environment=ALLOWED_ORIGIN=same-origin

WorkingDirectory=/opt/cortexos/packages/cortex-agent-generator
ExecStart=/usr/bin/node ./src/server.js

Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

NoNewPrivileges=no
ProtectSystem=false
ProtectHome=read-only
PrivateTmp=yes
RestrictRealtime=yes
LockPersonality=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
