# templates/systemd/9router.service
[Unit]
Description=9Router AI Gateway
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User={VPS_USER}
WorkingDirectory={VPS_HOME}
# Bind loopback only; Tailscale Serve is the public entrypoint. Upstream npm
# 9Router defaults to 0.0.0.0:20128, so keep CortexOS's canonical listener
# explicit here. Headless flags are mandatory for systemd.
ExecStart={NPM_PREFIX}/bin/9router --host 127.0.0.1 --port 11434 --skip-update --no-browser
Restart=always
RestartSec=5
Environment=HOME={VPS_HOME}
Environment=PATH={NODE_BIN_DIR}:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin
# CortexOS integration env. The systemd ExecStart above is the source of truth
# for host/port; NINEROUTER_BASE_URL remains the downstream client contract.
EnvironmentFile=-/opt/cortexos/.secrets/9router.env

[Install]
WantedBy=multi-user.target
