# templates/systemd/9router-docker-proxy.service
[Unit]
Description=9Router Docker bridge proxy
After=network-online.target 9router.service docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/cortexos/scripts/9router-docker-proxy.mjs
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
