# templates/systemd/9router-cloudflared.service
[Unit]
Description=Cloudflare quick tunnel for 9Router (port 11434)
After=network-online.target 9router.service
Wants=network-online.target
Requires=9router.service

[Service]
Type=simple
User=root
Group=root
Restart=on-failure
RestartSec=10
ExecStart=/opt/cortexos/.9router/bin/cloudflared tunnel --url http://127.0.0.1:11434 --no-autoupdate
StandardOutput=journal
StandardError=journal
SyslogIdentifier=9router-cloudflared

[Install]
WantedBy=multi-user.target
