# templates/systemd/cortex-degraded-service-watcher.service
[Unit]
Description=CortexOS degraded service AI watcher
Wants=network-online.target docker.service 9router.service
After=network-online.target docker.service 9router.service

[Service]
Type=oneshot
User=root
EnvironmentFile=-/opt/cortexos/.secrets/9router.env
EnvironmentFile=-/opt/cortexos/.secrets/dashboard.env
EnvironmentFile=-/opt/cortexos/.secrets/paperclip-app.env
WorkingDirectory=/opt/cortexos
ExecStart=/usr/bin/env node /opt/cortexos/scripts/cortex-degraded-service-watcher.mjs
