# templates/systemd/cortex-backup.service
[Unit]
Description=CortexOS encrypted full backup
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/cortexos/scripts/cortex-backup.sh
