# templates/systemd/cortex-mail-guardian-sweep.service
[Unit]
Description=Cortex Mail Guardian reconciliation sweep
After=network-online.target postgresql.service 9router.service
Wants=network-online.target

[Service]
Type=oneshot
User={VPS_USER}
Group={VPS_USER}
EnvironmentFile={CORTEX_SECRETS_DIR}/9router.env
EnvironmentFile={CORTEX_SECRETS_DIR}/dashboard.env
EnvironmentFile=-{CORTEX_SECRETS_DIR}/hermes/default.env
EnvironmentFile={CORTEX_SECRETS_DIR}/mail-guardian.env
WorkingDirectory={CORTEX_RUNTIME_ROOT}/packages/cortex-mail-guardian
ExecStart=/usr/bin/env node {CORTEX_RUNTIME_ROOT}/packages/cortex-mail-guardian/dist/index.js sweep
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths={CORTEX_SECRETS_DIR}
