# templates/systemd/cortex-auto-update.service
[Unit]
Description=CortexOS automatic package and service update
Wants=network-online.target docker.service
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/opt/cortexos/scripts/ops/cortex-auto-update.sh
