# templates/systemd/cortex-9remote.service
# 9Remote — remote terminal access (https://9remote.cc) systemd unit.
#
# Rendered by scripts/ops/cortex-render-units.sh which substitutes
# {CORTEX_ROOT} / {NPM_PREFIX} at deploy time. After rendering, reload systemd:
#
#   sudo bash scripts/ops/cortex-render-units.sh cortex-9remote.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now cortex-9remote.service
#
# Runs `9remote start` — the non-interactive server mode. Unlike the bare
# `9remote` TUI, `start` mode serves the web UI and spins up a Cloudflare
# tunnel for remote access. Use `9remote ui` instead if you only want
# tailnet-local access without the tunnel.
#
# Runs as the unprivileged `cortexos` user. This tool hands out interactive
# shells by design, so it is intentionally NOT sandboxed (ProtectHome/
# ProtectSystem would defeat the point). HOME points at a dedicated dir so
# the machine-bound permanent key persists across restarts.

[Unit]
Description=CortexOS 9remote browser terminal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={VPS_USER}
Group={VPS_USER}
Environment=HOME={VPS_HOME}
Environment=PORT=2208
ExecStart={NPM_PREFIX}/bin/9remote start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
