# templates/systemd/cortex-dashboard.service
[Unit]
Description=CortexOS Dashboard (native Next.js control plane)
Documentation=file://{CORTEX_ROOT}/packages/cortex-dashboard/AGENTS.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# root: the dashboard enters host namespaces directly (nsenter --target 1 in
# host-exec.ts), talks to the docker socket, and the root-helper socket. The
# former container ran privileged+pid:host as root; native root is equivalent
# and drops the privileged-container layer.
User=root
WorkingDirectory={CORTEX_ROOT}/packages/cortex-dashboard
EnvironmentFile={CORTEX_SECRETS_DIR}/dashboard.env
Environment=NODE_ENV=production
Environment=PORT=3080
# 0.0.0.0 so cortex-net producers (sandbox runner) reach the audit sink; the
# API is session-gated and the audit endpoint is Bearer-token gated.
Environment=HOSTNAME=0.0.0.0
Environment=DB_HOST=127.0.0.1
Environment=SETUP_STATE_PATH=/run/cortexos/setup-state.json
Environment=NEXT_TELEMETRY_DISABLED=1
# Migrations + spoke-aware catalog seed run before the server (mirrors the old
# docker-entrypoint.sh). migrate is fatal; dynamic-seed is best-effort.
ExecStartPre=/usr/bin/node scripts/migrate.js
ExecStartPre=-/usr/bin/node scripts/dynamic-seed.js
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
# Restart=always also rides out a not-yet-ready Postgres at boot.

[Install]
WantedBy=multi-user.target
