# templates/systemd/cortex-update-check.timer
[Unit]
Description=Run CortexOS native app update check daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
