# templates/systemd/cortex-mail-guardian-sweep.service
[Unit]
Description=Cortex Mail Guardian reconciliation sweep
After=network-online.target postgresql.service 9router.service
Wants=network-online.target

[Service]
Type=oneshot
User=cortexos
Group=cortexos
EnvironmentFile=/opt/cortexos/.secrets/9router.env
EnvironmentFile=/opt/cortexos/.secrets/dashboard.env
EnvironmentFile=-/opt/cortexos/.secrets/hermes/default.env
EnvironmentFile=/opt/cortexos/.secrets/mail-guardian.env
WorkingDirectory=/opt/cortexos/packages/cortex-mail-guardian
ExecStart=/usr/bin/env node /opt/cortexos/packages/cortex-mail-guardian/dist/index.js sweep
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/opt/cortexos/.secrets
