# templates/systemd/cortex-update-check.service
[Unit]
Description=CortexOS native app update check

[Service]
Type=oneshot
ExecStart=/opt/cortexos/scripts/cortex-update-check.sh
