# templates/systemd/cortex-mail-guardian-sweep.timer
[Unit]
Description=Run Cortex Mail Guardian sweep every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
