# templates/systemd/cortex-dashboard.service
[Unit]
Description=Cortex Dashboard
After=network-online.target postgresql.service docker.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User={VPS_USER}
WorkingDirectory=/opt/cortexos/dashboard
ExecStart={NODE_BIN} server.js
Restart=on-failure
RestartSec=5
Environment=PORT=3080
Environment=NODE_ENV=production
Environment=PATH={NODE_BIN_DIR}:/usr/local/bin:/usr/bin:/bin
EnvironmentFile=/opt/cortexos/.secrets/dashboard.env

[Install]
WantedBy=multi-user.target
