# templates/systemd/hermes-dashboard.service
[Unit]
Description=CortexOS Hermes Web UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={VPS_USER}
Group={VPS_USER}
Environment=HOME=/home/{VPS_USER}
Environment=PATH=/home/{VPS_USER}/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/snap/bin
Environment=HERMES_PROFILE=cortex
Environment=HERMES_HOME={CORTEX_HERMES_ROOT}/profiles/cortex
WorkingDirectory={CORTEX_HERMES_ROOT}/profiles/cortex
ExecStart=/home/{VPS_USER}/.local/bin/hermes dashboard --host localhost --port 9119 --no-open --skip-build
ExecStop=/home/{VPS_USER}/.local/bin/hermes dashboard --stop
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths={CORTEX_HERMES_ROOT} /home/{VPS_USER}/.hermes /tmp

[Install]
WantedBy=multi-user.target
