# templates/systemd/hermes-webui.service
# Hermes Web UI — nesquena/hermes-webui systemd unit.
#
# Rendered by scripts/ops/cortex-render-units.sh which substitutes
# {CORTEX_ROOT} and {CORTEX_SECRETS_DIR} at deploy time. After rendering,
# reload systemd:
#
#   sudo bash scripts/ops/cortex-render-units.sh hermes-webui.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hermes-webui.service
#
# Runs `docker compose up -d` from {CORTEX_ROOT}/hermes-webui (the install
# dir written by prompts/tools/30-hermes-webui.md). The compose file
# bind-mounts /opt/cortexos/.secrets/hermes-webui.env so HERMES_WEBUI_PASSWORD
# is never checked into git.

[Unit]
Description=Hermes Web UI (nesquena/hermes-webui)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=hermes-webui
Group=hermes-webui
WorkingDirectory={CORTEX_ROOT}/hermes-webui
EnvironmentFile={CORTEX_SECRETS_DIR}/hermes-webui.env
ExecStart=/usr/bin/docker compose up -d --remove-orphans --wait
ExecStop=/usr/bin/docker compose down
ExecReload=/bin/sh -c '/usr/bin/docker compose pull && /usr/bin/docker compose up -d --remove-orphans'
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
