# templates/systemd/ollama-honcho-embeddings-proxy.service
[Unit]
Description=Ollama embeddings proxy for Honcho Docker network
After=network-online.target docker.service ollama.service
Wants=network-online.target
Requires=ollama.service

[Service]
Type=simple
ExecStart=/usr/bin/socat TCP-LISTEN:11435,bind=172.30.0.1,reuseaddr,fork TCP:127.0.0.1:11435
Restart=always
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
