# templates/systemd/openviking-console.service
[Unit]
Description=OpenViking web console
After=network-online.target openviking.service
Wants=network-online.target openviking.service

[Service]
Type=simple
User=cortex
Group=cortex
WorkingDirectory=/opt/cortexos/stacks/openviking
EnvironmentFile=/opt/cortexos/.secrets/openviking.env
ExecStart=/opt/cortexos/stacks/openviking/.venv/bin/python -m openviking.console.bootstrap --host 127.0.0.1 --port 8020 --openviking-url http://127.0.0.1:18790
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
