# templates/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway
After=network-online.target docker.service
Wants=network-online.target docker.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User={VPS_USER}
Group={VPS_USER}
Environment=HOME={VPS_HOME}
Environment=OPENCLAW_HOME={VPS_HOME}
Environment=PATH={NODE_BIN_DIR}:/home/linuxbrew/.linuxbrew/bin:{VPS_HOME}/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=OPENVIKING_URL=http://127.0.0.1:18790
Environment=NINEROUTER_BASE_URL=http://127.0.0.1:11434
EnvironmentFile=-/opt/cortexos/.secrets/openclaw-gateway.env
ExecStart={NODE_BIN} {NPM_PREFIX}/lib/node_modules/openclaw/dist/index.js gateway --port 18789
ExecStop={NODE_BIN} {NPM_PREFIX}/lib/node_modules/openclaw/dist/index.js gateway stop
Restart=on-failure
RestartSec=10s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
