# templates/systemd/paperclip-agent-url-migrate.service
[Unit]
Description=Paperclip agent paperclipApiUrl auto-migration (admin-only, manual run)
Documentation=file:/opt/cortexos/paperclip/instances/default/RUNBOOK-paperclip-api-url.md
After=network-online.target paperclip.service
Wants=paperclip.service

[Service]
Type=oneshot
EnvironmentFile=/opt/cortexos/.secrets/paperclip.env
# Default: dry-run for safety. Override with `systemctl set-environment` or
# call the script directly to actually patch:
#   sudo PAPERCLIP_API_KEY=... node /opt/cortexos/scripts/paperclip-agent-url-migrate.mjs
Environment=PAPERCLIP_MIGRATE_DRY_RUN=1
ExecStart=/usr/bin/env node /opt/cortexos/scripts/paperclip-agent-url-migrate.mjs
StandardOutput=journal
StandardError=journal
