# templates/systemd/paperclip-local-proxy.socket
[Unit]
Description=Paperclip local compatibility proxy socket

[Socket]
ListenStream=127.0.0.1:3033
NoDelay=true

[Install]
WantedBy=sockets.target
