# stacks/cortex-agentgateway/cortex-agentgateway.service
[Unit]
Description=Cortex AgentGateway
After=network-online.target nats.service postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cortexos/stacks/cortex-agentgateway
EnvironmentFile=/opt/cortexos/.secrets/agentgateway.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/env node index.js
Restart=on-failure
RestartSec=5
User=cortexos
Group=cortexos
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
