# stacks/cortex-dashboard-root-helper/cortex-dashboard-root-helper.socket
[Unit]
Description=Cortex Dashboard Root Helper Socket

[Socket]
ListenStream=/run/cortexos/dashboard-helper.sock
SocketMode=0660
SocketUser=root
SocketGroup=cortexos
RemoveOnStop=yes
ExecStartPost=/usr/bin/chgrp cortexos /run/cortexos /run/cortexos/dashboard-helper.sock
ExecStartPost=/usr/bin/chmod 0750 /run/cortexos
ExecStartPost=/usr/bin/chmod 0660 /run/cortexos/dashboard-helper.sock

[Install]
WantedBy=sockets.target
